# bridge-service/systemd/ma352-bridge.service
[Unit]
Description=MA-352 RS-232 Bridge Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ma352-bridge
ExecStart=/usr/bin/python3 /opt/ma352-bridge/app.py
Restart=on-failure
RestartSec=2
Environment=BRIDGE_HOST=0.0.0.0
Environment=BRIDGE_PORT=5000
Environment=SERIAL_PORT=/dev/ttyUSB0
Environment=SERIAL_BAUD=115200
Environment=HOLD_INTERVAL=0.12
Environment=QUERY_INTERVAL=5.0
Environment=QUERY_ON_CONNECT=1

[Install]
WantedBy=multi-user.target
