# bridge-service/systemd/ma352-bridge.service
[Unit]
Description=MA-352 RS-232 Bridge Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ma352-bridge
ExecStart=/opt/ma352-bridge/.venv/bin/python /opt/ma352-bridge/app.py
Restart=on-failure
RestartSec=2
EnvironmentFile=-/etc/default/ma352-bridge

[Install]
WantedBy=multi-user.target
